# configs/essence.conf
[program:essence]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.py
directory=/var/www/essence/configs
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/essence.log